# .docker/php/crontab
0 8 * * * cd /var/www && php bin/console billy:commande:symfony
